# service/tailscale/run
#!/bin/sh
set -e
exec tailscaled -tun userspace-networking
